# .config/systemd/user/backup.service
[Unit]
Description=Backup configurations and files

[Service]
ExecStart=/home/eren/.config/systemd/user/backup.sh

[Install]
WantedBy=default.target
