# .config/systemd/user/backup.timer
[Unit]
Description=Backup files

[Timer]
OnCalendar=*-*-* *:00/5:00
Unit=backup.service

[Install]
WantedBy=default.target
