# .config/systemd/user/backup.path
[Unit]
Description=Monitor backup directories

[Path]
PathChanged=/home/eren/ssd/drive-cipher
PathChanged=/home/eren/ssd/media-cipher
PathChanged=/home/eren/ssd/mail-cipher
PathChanged=/home/eren/ssd/wallpapers
TriggerLimitIntervalSec=30s
Unit=backup.service

[Install]
WantedBy=default.target
