# .config/systemd/user/backup.timer
[Unit]
Description=Backup files

[Timer]
OnCalendar=minutely
Unit=backup.service

[Install]
WantedBy=default.target
