# .config/systemd/user/baser.path
[Unit]
Description=Monitor baser

[Path]
PathChanged=/home/eren/.dotfiles/arch/baser/
Unit=baser.service

[Install]
WantedBy=default.target
