# .config/systemd/user/baser.service
[Unit]
Description=Build baser

[Service]
ExecStart=%h/arch/baser/build.sh
WorkingDirectory=%h/arch/baser 

[Install]
WantedBy=default.target
