# .config/systemd/user/dotfiles.timer
[Unit]
Description=Backup configurations

[Timer]
OnCalendar=daily
Unit=dotfiles.service

[Install]
WantedBy=default.target
