# .config/systemd/user/dotfiles.service
[Unit]
Description=Backup configurations

[Service]
ExecStart= /bin/zsh -i -c "gitcom autosave /home/eren"

[Install]
WantedBy=default.target
