# .config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Session
Requires=basic.target
BindsTo=graphical-session.target
