# .config/systemd/user/launcher.path
[Unit]
Description=Monitor application directories

[Path]
PathChanged=/usr/share/applications
PathChanged=/home/eren/.local/share/applications
PathChanged=/home/eren/.local/bin/launcher
Unit=launcher.service

[Install]
WantedBy=default.target
