# .config/systemd/user/wallpaper.service
[Unit]
Description=Change the wallpaper

[Service]
ExecStart=/home/eren/.local/bin/wallpaper

[Install]
WantedBy=default.target
