# .config/autostart/.xdual.sh_1_.desktop
[Desktop Entry]
Exec=/home/eren/.xdual.sh
Icon=dialog-scripts
Name=.xdual.sh_1_
Type=Application
X-KDE-AutostartScript=true
