# .config/autostart/start.sh.desktop
[Desktop Entry]
Exec=/home/eren/.config/polybar/start.sh
Icon=dialog-scripts
Name=start.sh
Type=Application
X-KDE-AutostartScript=true
